# etc/rc.local
udevadm hwdb --update &
/usr/local/bin/heads-generate-passphrase
/usr/local/bin/heads-sup-compile &

exit 0
