# etc/rc.local
#[ -f /etc/unbound/root.key ] || {
#	/etc/init.d/unbound stop
#	unbound-anchor -a /etc/unbound/root.key
#	/etc/init.d/unbound start
#	return 0;
#}

#/usr/local/bin/heads-generate-passphrase

#/usr/local/bin/heads-torify -l &

exit 0
