# etc/rc.local
#!/bin/sh

/bin/udevadm hwdb --update &
/usr/local/lib/heads/heads-generate-passphrase
/usr/local/lib/heads/heads-lkm-compile 2>&1 >/dev/null &
/usr/local/lib/heads/heads-sup-compile 2>&1 >/dev/null &

for i in $(cat /proc/cmdline); do
	[ "$i" = "awesome" ] && wm=awesome
done

case "$wm" in
	awesome)
		cp /usr/share/heads/awesome.xinitrc /home/luther/.xinitrc
		cp /usr/share/heads/welcome-awesome.txt /usr/share/heads/welcome.txt
		;;
	*)
		cp /usr/share/heads/openbox.xinitrc /home/luther/.xinitrc
		cp /usr/share/heads/welcome-openbox.txt /usr/share/heads/welcome.txt
		;;
esac

exit 0
