# etc/rc.local
udevadm hwdb --update &
/usr/local/bin/heads-generate-passphrase

exit 0
