# etc/rc.local
[ -f /etc/unbound/root.key ] || {
	#/etc/init.d/unbound stop
	#unbound-anchor -a /etc/unbound/root.key
	#/etc/init.d/unbound start
	return 0;
}
