# etc/rc.local
udevadm hwdb --update &
/usr/local/lib/heads/heads-generate-passphrase
/usr/local/lib/heads/heads-sup-compile &

exit 0
